# deploy/claude-bot.service
[Unit]
Description=Claude AI Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/test-portfolio
EnvironmentFile=/root/test-portfolio/.env
ExecStart=/root/test-portfolio/.venv/bin/python -m bot.main
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
